# less/templates/upstart.conf
description "Monitors for .less changes"
start on startup
respawn
task
exec <%= @script %>
